# src/cloudcrypt/etc/cloudcrypt.service
[Unit]
Description=Simple Cloud Encryption Service

[Service]
Type=exec
ExecStart=$ScriptPath

[Install]
WantedBy=multi-user.target